# deploy/ngrok.service
[Unit]
Description=ngrok tunnel for the SMS webhook (static domain)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/printerservice
ExecStart=/usr/local/bin/ngrok start --all --config /home/admin/printerservice/ngrok.yml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
