# deploy/sms-listener.service
[Unit]
Description=Twilio SMS webhook -> thermal printer (Express)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/printerservice
# Node is installed via nvm, so there is no /usr/bin/node. Sourcing nvm.sh puts
# the default node version on PATH (survives nvm upgrades); exec replaces the
# shell so systemd signals reach node directly.
ExecStart=/bin/bash -c '. /home/admin/.nvm/nvm.sh && exec node src/server.js'
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
