# deploy/sms-listener.service
[Unit]
Description=Twilio SMS webhook -> thermal printer (Express)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/printerservice
ExecStart=/usr/bin/node src/server.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
